# scripts/drowsiness_detection.service
# scripts/drowsiness-detection.service
[Unit]
Description=Driver Drowsiness Detection System
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/driver-drowsiness-detection
Environment=PATH=/home/pi/driver-drowsiness-detection/drowsiness_env/bin
ExecStart=/home/pi/driver-drowsiness-detection/drowsiness_env/bin/python src/main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
